# s6-rc.d/dovecot/run
#!/command/execlineb -P

foreground { chown -R mail:mail u=rwX,g=rwX,o= /app/mail }
# this is necessary because otherwise dovecot fails with some weird "failed to dump capabilities error"
foreground { sleep 1 }

dovecot -F
