# services.d/redis_rspamd/run
#!/usr/bin/execlineb -P
redis-server /etc/rspamd/redis_server.conf
