# s6-rc.d/opendkim/run
#!/command/execlineb -P

foreground { ./data/copy_opendkim_conf.sh }

opendkim -f -x /etc/opendkim/opendkim.conf
