# services.d/postfix/run
#!/usr/bin/execlineb -P
postfix start-fg
