# s6-rc.d/rspamd_redis/run
#!/bin/bash
set -e

ln -s /app/conf/rspamd_worker_controller.inc /etc/rspamd/override.d/worker-controller.inc
mkdir -p /app/data/rspamd/redis

exec redis-server /etc/rspamd/redis_server.conf
