# services.d/opendkim/run
#!/usr/bin/execlineb -P
opendkim -f -x /etc/opendkim/opendkim.conf
